# vypr-luci/files/rc.local
# Put your custom commands here that should be executed once
# the system init finished. By default this file does nothing.

# /etc/init.d/openvpn enable
/usr/bin/check_openvpn.sh 30 &

exit 0